# docker/images/yt-dlp/crontab
1 */4 * * * python3 /usr/bin/kinocheck-trailers.py
